# docker/build/akeneo_daemon.conf
[program:akeneo_queue_daemon]
command=/usr/bin/php /srv/pim/bin/console akeneo:batch:job-queue-consumer-daemon --env=prod
autostart=false
autorestart=true
stderr_logfile=/var/log/akeneo_daemon.err.log
stdout_logfile=/var/log/akeneo_daemon.out.log
user=www-data
